# rc.local
#!/bin/sh -e
#
# rc.local

python3 /home/root/main_loop_config30bins.py >> /home/root/main.log 2>&1 &
python3 /home/root/button.py >> /home/root/button.log 2>&1 &

exit 0

